# managed/devops/roles/configure-cluster-server/templates/yb-clean_cores.service
[Unit]
Description=Yugabyte clean cores

[Service]
{% if ansible_os_family == 'RedHat' and (ansible_distribution_major_version == '7' or (ansible_distribution == 'Amazon' and ansible_distribution_major_version == '2')) %}
User={{ user_name }}
Group={{ user_name }}
{% endif %}
Type=oneshot
WorkingDirectory={{yb_home_dir}}/bin
ExecStart=/bin/bash {{yb_home_dir}}/bin/clean_cores.sh

[Install]
WantedBy=multi-user.target
